# scenarios/foo-bar-realtime-streaming-metrics/supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:redpanda]
command=rpk redpanda start --mode dev-container
autostart=true
autorestart=false
priority=10
stdout_logfile=/tmp/redpanda.log
stderr_logfile=/tmp/redpanda.err.log

[program:clickhouse]
command=/usr/bin/clickhouse-server --config-file=/etc/clickhouse-server/config.xml
user=clickhouse
autostart=true
autorestart=false
priority=20
stdout_logfile=/tmp/clickhouse.log
stderr_logfile=/tmp/clickhouse.err.log
